# etc/nodepool-launcher.service.d/centos.conf
[Service]
ExecStart=
ExecStart=/usr/bin/nodepool-launcher
